# server/etc/systemd/system/autousbip@.service
[Unit]
Description=autousbip trigger

[Service]
ExecStart=/usr/local/bin/autousbip.sh %i 192.168.0.5

[Install]
WantedBy=multi-user.target
